# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
user=root

[program:oracle]
command=/bin/bash -c "/oracle_installation/oracle.sh"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0